# docs/dfir/attacks/persesc/index.xml
<?xml version="1.0" encoding="utf-8" standalone="yes"?>
<rss version="2.0" xmlns:atom="http://www.w3.org/2005/Atom">
  <channel>
    <title>Persistence and Escalation Mechanisms on Analyst</title>
    <link>https://bakerst221b.com/docs/dfir/attacks/persesc/</link>
    <description>Recent content in Persistence and Escalation Mechanisms on Analyst</description>
    <generator>Hugo -- gohugo.io</generator>
    <language>en-us</language>
    <lastBuildDate>Fri, 26 May 2023 00:00:00 +0000</lastBuildDate><atom:link href="https://bakerst221b.com/docs/dfir/attacks/persesc/index.xml" rel="self" type="application/rss+xml" />
    <item>
      <title>Collection</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/collection/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/collection/</guid>
      <description>&lt;p&gt;*There are a lot of techniques to become persistent on a Windows machine. You can refer to the persistence and escalation mechanisms article or look those up on the MITRE website. Although it&amp;rsquo;s important to know the artefacts themselves and the technique, it would be hard to collect them all manually. *&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>BIOS Persistence</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/bios/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/bios/</guid>
      <description>&lt;p&gt;&lt;em&gt;This is about &amp;hellip; .&lt;/em&gt;&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Boot Abuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/boot/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/boot/</guid>
      <description>Linux Init is the process that spawns other processes and usually has pid 1. /etc/inittab can be used to set run level for the systems.</description>
    </item>
    
    <item>
      <title>Config Abuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/startup/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/startup/</guid>
      <description>&lt;p&gt;&lt;em&gt;Some persistence and escalation techniques don&amp;rsquo;t require an exploit but only a little tweaking of some config.&lt;/em&gt;&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Elevation Control Abuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/elevation/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/elevation/</guid>
      <description>Linux &amp;amp; macOS setuid and setgid Platforms: macOS, Linux MITRE: https://attack.mitre.org/techniques/T1548/001/
setuid or setgid bits set in UNIX.</description>
    </item>
    
    <item>
      <title>IPC</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/ipc/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/ipc/</guid>
      <description>&lt;p&gt;&lt;em&gt;This is about &amp;hellip; .&lt;/em&gt;&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Kernel</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/kernel/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/kernel/</guid>
      <description>Linux LKM LKM for Linux. Drivers are one type of kernel extentions.
XDG https://attack.mitre.org/techniques/T1547/013/
macOS kext kext for macOS.</description>
    </item>
    
    <item>
      <title>Library Abuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/libs/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/libs/</guid>
      <description>Windows DLL Search Order Hijacking 👑 - require w permissions for Windows and Windows\System32.
Detection:</description>
    </item>
    
    <item>
      <title>Services</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/services/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/services/</guid>
      <description>&lt;p&gt;&lt;em&gt;Services can be abused differently but are often used as escalation and persistence mechanisms.&lt;/em&gt;&lt;/p&gt;</description>
    </item>
    
    <item>
      <title>Tasks Abuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/tasks/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/tasks/</guid>
      <description>Windows ❗️ Both at and schtasks.exe can create tasks remotely.
🏺 Artefacts:
🪵 Event Log: Task Scheduler Operational.</description>
    </item>
    
    <item>
      <title>Tokens Misuse</title>
      <link>https://bakerst221b.com/docs/dfir/attacks/persesc/tokens/</link>
      <pubDate>Sat, 03 Jun 2023 00:00:00 +0000</pubDate>
      
      <guid>https://bakerst221b.com/docs/dfir/attacks/persesc/tokens/</guid>
      <description>MITRE: https://attack.mitre.org/techniques/T1134/001/, https://attack.mitre.org/techniques/T1134/, https://attack.mitre.org/techniques/T1134/002/, https://attack.mitre.org/techniques/T1134/003/, https://attack.mitre.org/techniques/T1134/004/, https://attack.mitre.org/techniques/T1134/005/ Actors: https://attack.mitre.org/groups/G0032/
Techniques:
Token impersonalisation. runas or CreateProcessWithTokenW to create a process with the rights of another user.</description>
    </item>
    
  </channel>
</rss>
